# collector/maiwar_wq.service
[Unit]
Description=Fetch new measurements data for Maiwar WQ

[Service]
ExecStart=maiwar_wq_collector
