# collector/maiwar_wq.service
[Unit]
Description=Fetch new measurements data for Maiwar WQ

[Service]
ExecStart=/path/to/maiwar_wq_collector
